# scripts/pos-backend.service
[Unit]
Description=Restaurant POS Backend Service
After=network.target

[Service]
User=pos
Group=pos
WorkingDirectory=/home/pos/AppPOS
Environment=PATH=/home/pos/AppPOS/venv/bin:$PATH
ExecStart=/home/pos/AppPOS/venv/bin/python -m uvicorn app:app --host 0.0.0.0 --port 8000 --workers 1
Restart=always

[Install]
WantedBy=multi-user.target 